# rc.d/gohttpd
#!/bin/ksh
#
# OpenBSD rc.d file for go-httpd
# put into /etc/rc.d

daemon="/usr/local/sbin/go-httpd"

. /etc/rc.d/rc.subr

rc_reload=NO
rc_bg=YES

rc_cmd $1
